# rc.local
/usr/local/bin/uwsgi --ini /APPDIR/APPNAME.ini --uid www-data --gid www-data --daemonize /var/log/uwsgi.log

exit 0
